# JAEGIS\JAEGIS-METHOD\emad-auto-sync.service
[Unit]
Description=eJAEGIS Auto-Sync Monitoring Service
Documentation=https://github.com/YOUR_USERNAME/eJAEGIS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=eJAEGIS
Group=eJAEGIS
WorkingDirectory=/path/to/JAEGIS-METHOD
ExecStart=/usr/bin/python3 /path/to/JAEGIS-METHOD/eJAEGIS-auto-sync.py --jaegis-path /path/to/JAEGIS-METHOD --interval 3600
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eJAEGIS-auto-sync

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/JAEGIS-METHOD
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=PYTHONPATH=/path/to/JAEGIS-METHOD
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
